# deploy/boilerhouse-podmand.service
[Unit]
Description=Boilerhouse Runtime Daemon (podman + privilege-separating proxy)
Documentation=https://docs.podman.io/en/latest/markdown/podman-system-service.1.html
After=network.target

[Service]
Type=simple
UMask=0117
Environment=PODMAN_SOCKET=/run/boilerhouse/podman.sock
Environment=LISTEN_SOCKET=/run/boilerhouse/runtime.sock
Environment=SNAPSHOT_DIR=/var/lib/boilerhouse/snapshots
ExecStartPre=/bin/mkdir -p /run/boilerhouse /var/lib/boilerhouse/snapshots
ExecStart=/usr/local/bin/bun /opt/boilerhouse/apps/boilerhouse-podmand/src/main.ts
ExecStartPost=/bin/sh -c 'until [ -S /run/boilerhouse/runtime.sock ]; do sleep 0.1; done; chgrp %i /run/boilerhouse/runtime.sock; chmod 660 /run/boilerhouse/runtime.sock'
Restart=on-failure
RestartSec=5

# Security hardening
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_SYS_PTRACE CAP_NET_ADMIN CAP_SYS_CHROOT CAP_CHOWN CAP_DAC_READ_SEARCH CAP_FOWNER CAP_KILL
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/run/boilerhouse /var/lib/boilerhouse /var/lib/containers /var/run/containers /tmp

[Install]
WantedBy=multi-user.target
