# deploy/boilerhouse-api.service
[Unit]
Description=Boilerhouse API
After=network-online.target boilerhouse-podmand.service
Wants=network-online.target
Requires=boilerhouse-podmand.service

[Service]
Type=simple
User=boilerhouse
Group=boilerhouse
WorkingDirectory=/opt/boilerhouse
ExecStart=/usr/local/bin/bun apps/api/src/server.ts
Restart=on-failure
RestartSec=5

# Environment
EnvironmentFile=/etc/boilerhouse/api.env

# Hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
MemoryDenyWriteExecute=no
ReadWritePaths=/opt/boilerhouse/data

[Install]
WantedBy=multi-user.target
